# api/iobsp-coach-api.service
[Unit]
Description=IOBSP Coach API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/iobsp-coach-api
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/opt/iobsp-coach-api/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
